# ibkr-api.service
[Unit]
Description=IBKR AI Trader API + Scheduler
After=network.target
Wants=network.target

[Service]
Type=simple
User=frankpach
WorkingDirectory=/home/frankpach/ibkr-bot
EnvironmentFile=/home/frankpach/ibkr-bot/.env
ExecStart=/home/frankpach/ibkr-bot/.venv/bin/python3 /home/frankpach/ibkr-bot/run.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
