# ibkr-gateway.service
[Unit]
Description=Interactive Brokers Gateway
After=network.target
Wants=network.target

[Service]
Type=simple
User=frankpach
Environment=DISPLAY=:0
ExecStart=/home/frankpach/Jts/ibgateway/1046/ibgateway
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
